# deploy/swiftpropdf-worker.service
[Unit]
Description=SwiftProPDF Celery Worker
Wants=network-online.target
After=network-online.target postgresql.service redis-server.service clamav-daemon.service

[Service]
User=ubuntu
Group=ubuntu

WorkingDirectory=/home/ubuntu/SwiftProPDF
EnvironmentFile=/home/ubuntu/SwiftProPDF/.env
Environment=PYTHONUNBUFFERED=1

ExecStart=/bin/sh -c 'exec /home/ubuntu/venv/bin/celery -A SwiftProPDF.web_app.celery_app:celery_app worker --loglevel=info --concurrency="${SWIFTPROPDF_ASYNC_WORKER_CONCURRENCY:-2}"'

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
